# Ansible/services/crondiff.timer
[Unit]
Description=Crondiff -- Timer part

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
