# Ansible/services/crondiff.service
[Unit]
Description=Mails root if there's a diff in crontab

[Service]
User=root
ExecStart=/opt/crondiff
Type=oneshot
